# docker/php/crontab
# Crontab
# This crontab runs inside of the PHP container

* * * * * cd /var/www && php ./vendor/bin/crunz schedule:run
